# ipcam-flask.service
[Unit]
Description=IPCam PWA — Flask server
After=network.target ipcam-ffmpeg.service

[Service]
Type=simple
User=mdevuono_pi
WorkingDirectory=/home/mdevuono_pi/ipcam-pwa
ExecStart=/home/mdevuono_pi/ipcam-pwa/venv/bin/python app.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
