# ipcam-pwa/ipcam-ffmpeg.service
[Unit]
Description=IPCam ffmpeg HLS encoder
After=network.target

[Service]
Type=simple
User=mdevuono_pi
ExecStartPre=/bin/mkdir -p /tmp/hls
ExecStartPre=/bin/bash -c 'rm -f /tmp/hls/*.ts /tmp/hls/*.m3u8'
ExecStart=/usr/bin/ffmpeg \
  -f v4l2 \
  -input_format mjpeg \
  -video_size 1920x1080 \
  -framerate 30 \
  -i /dev/video0 \
  -vf format=yuv420p \
  -c:v libx264 \
  -preset ultrafast \
  -tune zerolatency \
  -profile:v baseline \
  -level 3.1 \
  -b:v 4M \
  -g 60 \
  -keyint_min 60 \
  -sc_threshold 0 \
  -bsf:v h264_mp4toannexb \
  -hls_time 2 \
  -hls_list_size 3 \
  -hls_flags delete_segments+omit_endlist \
  -hls_segment_filename /tmp/hls/seg%%03d.ts \
  -f hls \
  /tmp/hls/stream.m3u8
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
